# nats.supervisor.conf
[program:nats]
command=/bin/bash -c "cd /app/nats && nice -n 5 ./gnatsd -c seed.conf"
redirect_stderr=true
